# LV2/pipitch-connect.service
[Unit]
Description=PiPitch MIDI fan-out — connect ZynMidiRouter:ch0_out to extra synth chains
After=zynthian.service
Requires=zynthian.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pipitch-connect.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
